# llama/llama-embed.service
[Unit]
Description=llama-server (embedding model)
After=default.target

[Service]
Type=simple
ExecStart=@SERVER_BIN@ \
    --model @MODEL_PATH@ \
    --port @PORT@ \
    --host 0.0.0.0 \
    --embeddings \
    --pooling last \
    --n-gpu-layers 0 \
    --ctx-size 2048
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
