# llama/llama-server.service
[Unit]
Description=llama-server (GPU-accelerated LLM inference)
After=default.target

[Service]
Type=simple
# Listen on Tailscale so laptop clients can use the desktop GPU server.
ExecStart=@SERVER_BIN@ \
    --model @MODEL_PATH@ \
    --port @PORT@ \
    --host 0.0.0.0 \
    --ctx-size 32768 \
    --n-gpu-layers 99 \
    --flash-attn on \
    --cache-type-k q8_0 \
    --cache-type-v q8_0 \
    --reasoning off
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
